# recipes-connectivity/nrf52-hci-uart/files/nrf52-btattach-hci-uart.service
[Unit]
Description=btattach hci_uart, needed to enable Bluetooth for nrf52832 zephyr hci uart adapter
After=bluetooth.target

[Service]
Type=simple
ExecStart=/usr/bin/btattach -B /dev/ttyHCIUART0  -S 1000000

[Install]
WantedBy=multi-user.target
