# recipes-support/nanopi-neo-air-fix/files/99-brcmfmac.rules
SUBSYSTEM=="net",  ACTION=="add", KERNEL=="wlan*", DRIVERS=="brcmfmac", RUN+="/usr/sbin/iw wlan0 set power_save off", NAME="wlan0"
